# systemd/catalyst-agent.service.d/limits.conf
[Service]
LimitNOFILE=1048576
LimitNOFILESoft=1048576
TasksMax=infinity
MemoryMax=4G
MemorySwapMax=0
CPUWeight=1000
IOWeight=1000
