# containerd/catalyst-mysql.service
[Unit]
Description=Catalyst MySQL Database (containerd)
After=network-online.target containerd.service
Requires=containerd.service
Wants=network-online.target

[Service]
Type=forking
Restart=always
RestartSec=10s
TimeoutStartSec=300
TimeoutStopSec=30

# Load secrets from env file (create /etc/catalyst/mysql.env with your credentials)
EnvironmentFile=-/etc/catalyst/mysql.env

# Environment variables (override in /etc/catalyst/mysql.env)
Environment="CONTAINER_NAME=catalyst-mysql"
Environment="MYSQL_IMAGE=docker.io/library/mysql:8.4"
Environment="MYSQL_PORT=3306"
Environment="MYSQL_DATABASE=catalyst_databases"
Environment="DATA_DIR=/var/lib/catalyst/mysql"
Environment="CTR_NAMESPACE=catalyst"
Environment="SCRIPT_DIR=/opt/catalyst/containerd"

# Pre-start: ensure directory exists
ExecStartPre=-/usr/bin/mkdir -p /var/lib/catalyst/mysql

# Main start command
ExecStart=${SCRIPT_DIR}/mysql-wrapper-start.sh

# Stop command
ExecStop=/usr/bin/ctr -n catalyst task kill catalyst-mysql SIGTERM
ExecStop=/bin/sleep 10
ExecStop=/usr/bin/ctr -n catalyst task rm -f catalyst-mysql || true

# Post-stop
ExecStopPost=/bin/rm -f /tmp/catalyst-mysql.pid

# Resource limits
LimitNOFILE=1048576
LimitNPROC=1048576

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
